# etc/systemd/user/nvidia-settings.service
[Unit]
ConditionPathExists=%h/.nvidia-settings-rc

[Service]
ExecStart=/opt/nvidia/nvidia-settings -l

[Install]
WantedBy=xsession.target
